# rootfs/etc/services.d/varnish/run
#!/usr/bin/with-contenv sh
exec /varnish.sh
