# rootfs/etc/services.d/prometheus_varnish_exporter/run
#!/usr/bin/execlineb -P
/usr/local/bin/prometheus_varnish_exporter -no-exit
