# rootfs/etc/services.d/prometheus_varnish_exporter/run
#!/usr/bin/with-contenv sh
: ${VARNISH_SKIP_METRICS:=false}
if [[ $VARNISH_SKIP_METRICS != "true" ]]; then
  exec /usr/local/bin/prometheus_varnish_exporter -no-exit
fi
