# rootfs/etc/services.d/prometheus_varnish_exporter/run
#!/usr/bin/with-contenv sh
exec /usr/local/bin/prometheus_varnish_exporter -no-exit
