# rootfs/etc/services.d/dns_check/run
#!/usr/bin/execlineb -P
exec /dnscheck.sh
