# rootfs/etc/services.d/dns_check/run
#!/usr/bin/with-contenv sh
/dnscheck.sh
