# rootfs/etc/services.d/dnscheck/run
#!/usr/bin/with-contenv sh
exec /dnscheck.sh
